# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:pboxa]
command=/usr/bin/pboxa run -config=/etc/pboxa/cf
autorestart=true
priority=200

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock
